# ftp/rootfs/etc/s6-overlay/s6-rc.d/vsftpd/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: FTP
# Runs vsftpd (IPv4)
# ==============================================================================
# Wait for vsftpd-log to be started
s6-svwait -u -t 5000 /run/service/vsftpd-log

bashio::log.info "Starting FTP server..."

# RUN vsftpd
exec vsftpd /etc/vsftpd/vsftpd.conf
