# ftp/rootfs/etc/services.d/vsftpd-log/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: FTP
# Catches output from vsFTPd server
# ==============================================================================
bashio::log.info "Starting Log tailer..."
touch /var/log/vsftpd.log
tail -n0 -F /var/log/vsftpd.log
