# ftp/rootfs/etc/services.d/vsftpd/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: FTP
# Runs vsftpd (IPv4)
# ==============================================================================
# Wait for vsftpd-log to be started
s6-svwait -u -t 5000 /var/run/s6/services/vsftpd-log

bashio::log.info "Starting FTP server..."

# RUN vsftpd
exec vsftpd /etc/vsftpd/vsftpd.conf
